# pgbackrest/pg_diff_backup.service
[Unit]
Description=PostgreSQL diff backup

[Service]
Type=simple
User=postgres
ExecStart=/etc/pgbackrest/scripts/pg_backup.sh --diff
