# pgbackrest/pg_full_backup.service
[Unit]
Description=PostgreSQL full backup

[Service]
Type=simple
User=postgres
ExecStart=/etc/pgbackrest/scripts/pg_backup.sh --full
